# kernel-update-hook.timer
[Unit]
Description=Update list of modules regularly.

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
