# kernel-update-hook.service
[Unit]
Description=Update list of used modules.

[Service]
ExecStart=/usr/bin/kernel-update-modules

[Install]
WantedBy=shutdown.target

